# oled-status@.service
[Unit]
Description=OLED status display (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/app
ExecStart=%h/app/.venv/bin/python -m src.oled_status --device i2c --addr 0x3C --refresh 5
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
