# rpi-gpio@.service
[Unit]
Description=RPi GPIO Web Interface (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# Use the service user's home, not /root
WorkingDirectory=%h/app

# Create the logs dir without using a shell or root-only chown
ExecStartPre=/usr/bin/mkdir -p %h/app/data/logs

# Start your app module
ExecStart=%h/app/.venv/bin/python -m src

Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8
Environment=PYTHONUNBUFFERED=1

Restart=on-failure
RestartSec=3
SupplementaryGroups=gpio

[Install]
WantedBy=multi-user.target
