# systemd/rpi-gpio@.service
[Unit]
Description=RPi GPIO Web Interface (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# %h == /home/<user>
WorkingDirectory=%h/app
Environment=LC_ALL=C.UTF-8
Environment=LANG=C.UTF-8
Environment=PYTHONUNBUFFERED=1

# Make logs dir owned by the service user (no $HOME, no /root)
ExecStartPre=/usr/bin/install -d -m 0755 -o %i -g %i %h/app/data/logs

# Adjust module if your entrypoint differs
ExecStart=%h/app/.venv/bin/python -m src

Restart=on-failure
RestartSec=3
SupplementaryGroups=gpio i2c

[Install]
WantedBy=multi-user.target
